# files/k3screenctrl.init
#!/bin/sh /etc/rc.common

START=99
STOP=99

TEMP_DIR=/tmp/k3screenctrl

start() {
    mkdir -p $TEMP_DIR
    mkdir -p $TEMP_DIR/device_speed
    /usr/bin/k3screenctrl &
}

stop() {
    killall k3screenctrl
}
